# systemd/dynamic-fan-ui.service
[Unit]
Description=Dynamic Fan Web UI (Flask)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
# For production, you might run via gunicorn; this simple unit uses Flask directly
WorkingDirectory=/opt/dynamic-fan-ui/app
ExecStart=/usr/bin/env python3 app.py
Restart=always
RestartSec=2
Environment=FAN_CURVE_FILE=/opt/dynamic-fan-ui/fan_curve.json
Environment=FAN_SPEED_FILE=/opt/dynamic-fan-ui/fan_speeds.txt
Environment=FAN_SPEED_BITS_FILE=/opt/dynamic-fan-ui/fan_speeds_bits.txt
# iLO config also used for temp reads
Environment=ILO_IP=192.168.1.100
Environment=ILO_USER=admin
Environment=ILO_SSH_KEY=/root/.ssh/ilo_key
# Optional password
# Environment=ILO_PASSWORD=
# Optional IPMI temps
# Environment=USE_IPMI_TEMPS=1
Environment=PWM_UNITS=bits

[Install]
WantedBy=multi-user.target
