# systemd/dynamic-fans.service
[Unit]
Description=Dynamic Fan Control for HP ML350p Gen8
After=network-online.target
Wants=network-online.target
Wants=dynamic-fan-ui.service

[Service]
Type=simple
ExecStart=/usr/local/bin/dynamic_fans.sh
Restart=always
RestartSec=2
# Environment overrides placed in /etc/systemd/system/dynamic-fans.service.d/override.conf
Environment=FAN_CURVE_FILE=/opt/dynamic-fan-ui/fan_curve.json
Environment=FAN_SPEED_FILE=/opt/dynamic-fan-ui/fan_speeds.txt
Environment=FAN_SPEED_BITS_FILE=/opt/dynamic-fan-ui/fan_speeds_bits.txt
Environment=PWM_UNITS=bits
Environment=CHECK_INTERVAL=1
Environment=MAX_STEP=20
Environment=ILO_CMD_GAP_MS=1000
# iLO auth defaults; override via override.conf for host-specific details
Environment=ILO_IP=192.168.1.100
Environment=ILO_USER=admin
Environment=ILO_SSH_KEY=/root/.ssh/ilo_key
# Optional password auth if sshpass available
# Environment=ILO_PASSWORD=

[Install]
WantedBy=multi-user.target
